# code/supervisord.conf
[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
loglevel=info
pidfile=/tmp/supervisord.pid

[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

# ─────────────────────────────────────────────
# FastAPI backend (uvicorn)
# Runs on internal port 8000
# ─────────────────────────────────────────────
[program:uvicorn]
command=python -m uvicorn main:app --host 127.0.0.1 --port 8000 --workers 1 --log-level info
directory=/app/backend
autostart=true
autorestart=true
startretries=5
startsecs=10
stopwaitsecs=30
stdout_logfile=/tmp/uvicorn.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile=/tmp/uvicorn_err.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
environment=PYTHONPATH="/app/backend",PYTHONUNBUFFERED="1"
priority=10

# ─────────────────────────────────────────────
# nginx (serves React + proxies /api/ to uvicorn)
# Listens on port 9000 (external)
# ─────────────────────────────────────────────
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
startretries=5
startsecs=5
stopwaitsecs=10
stdout_logfile=/tmp/nginx.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile=/tmp/nginx_err.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
priority=20

